# wazuh-ip-reputation.service
[Unit]
Description=Wazuh IP Reputation Checker
Documentation=https://github.com/wazuh/wazuh-ip-reputation
After=network.target network-online.target
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
Type=simple
User=wazuh-reputation
Group=wazuh-reputation
WorkingDirectory=/opt/wazuh-ip-reputation
Environment=PATH=/opt/wazuh-ip-reputation/venv/bin
Environment=PYTHONPATH=/opt/wazuh-ip-reputation
ExecStart=/opt/wazuh-ip-reputation/venv/bin/python3 /opt/wazuh-ip-reputation/wazuh_ip_reputation.py --continuous
ExecReload=/bin/kill -HUP $MAINPID
ExecStop=/bin/kill -TERM $MAINPID
Restart=always
RestartSec=10
TimeoutStopSec=30
KillMode=mixed
StandardOutput=journal
StandardError=journal
SyslogIdentifier=wazuh-ip-reputation

# Límites de recursos
LimitNOFILE=65536
LimitNPROC=32768

# Configuración de seguridad
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/wazuh-ip-reputation /var/log /tmp
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictNamespaces=true
RestrictSUIDSGID=true
RemoveIPC=true
PrivateDevices=true
ProtectHostname=true
ProtectClock=true
ProtectKernelModules=true
ProtectKernelLogs=true
MemoryDenyWriteExecute=true
LockPersonality=true
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources

# Configuración de red
IPAddressDeny=any
IPAddressAllow=localhost
IPAddressAllow=10.0.0.0/8
IPAddressAllow=172.16.0.0/12
IPAddressAllow=192.168.0.0/16

[Install]
WantedBy=multi-user.target
